# config/supervisor/backend.conf
[program:backend]
command=bash -c "sleep 5 && cd /app/backend && node dist/index.js"
user=root
autostart=true
autorestart=true
priority=300
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
directory=/app/backend
stopwaitsecs=30
environment=NODE_ENV="production",PORT="6002",DATABASE_URL="postgresql://postgres:transcribeai_postgres_pass@localhost:5432/transcribeai"
